# supervisor/conf.d/nginx.conf
[program:nginx]
command = /usr/sbin/nginx -g

autorestart = true

stderr_logfile = NONE
stdout_logfile = NONE