# supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command = php-fpm81 -F
autorestart = true